# tc-packages/virtualbox-ose-beta/extensions/virtualbox-ose-beta.list
lib/udev/rules.d/60-vboxguest.rules
lib/udev/rules.d/60-virtualbox.rules
usr/local/bin/VBoxClient
usr/local/bin/VBoxService
usr/local/bin/VBox
usr/local/bin/VBoxClient-all
usr/local/bin/VBoxControl
usr/local/bin/VirtualBox
usr/local/bin/VBoxHeadless
usr/local/bin/VBoxSDL
usr/local/bin/VBoxTunctl
usr/local/bin/VBoxManage
usr/local/share/doc/virtualbox-ose/license
usr/local/share/pixmaps/VBox.png
usr/local/share/virtualbox/nls/qt_zh_TW.qm
usr/local/share/virtualbox/nls/VirtualBox_hu.qm
usr/local/share/virtualbox/nls/qt_fi.qm
usr/local/share/virtualbox/nls/qt_ja.qm
usr/local/share/virtualbox/nls/qt_ca.qm
usr/local/share/virtualbox/nls/VirtualBox_ro.qm
usr/local/share/virtualbox/nls/qt_cs.qm
usr/local/share/virtualbox/nls/qt_km_KH.qm
usr/local/share/virtualbox/nls/VirtualBox_cs.qm
usr/local/share/virtualbox/nls/VirtualBox_es.qm
usr/local/share/virtualbox/nls/VirtualBox_sk.qm
usr/local/share/virtualbox/nls/VirtualBox_bg.qm
usr/local/share/virtualbox/nls/qt_de.qm
usr/local/share/virtualbox/nls/VirtualBox_tr.qm
usr/local/share/virtualbox/nls/VirtualBox_it.qm
usr/local/share/virtualbox/nls/qt_en.qm
usr/local/share/virtualbox/nls/qt_bg.qm
usr/local/share/virtualbox/nls/VirtualBox_ko.qm
usr/local/share/virtualbox/nls/VirtualBox_ca_VA.qm
usr/local/share/virtualbox/nls/qt_sr.qm
usr/local/share/virtualbox/nls/VirtualBox_fi.qm
usr/local/share/virtualbox/nls/VirtualBox_zh_CN.qm
usr/local/share/virtualbox/nls/qt_el.qm
usr/local/share/virtualbox/nls/VirtualBox_gl_ES.qm
usr/local/share/virtualbox/nls/qt_sk.qm
usr/local/share/virtualbox/nls/qt_nl.qm
usr/local/share/virtualbox/nls/VirtualBox_da.qm
usr/local/share/virtualbox/nls/VirtualBox_el.qm
usr/local/share/virtualbox/nls/qt_lt.qm
usr/local/share/virtualbox/nls/qt_ru.qm
usr/local/share/virtualbox/nls/qt_zh_CN.qm
usr/local/share/virtualbox/nls/qt_hu.qm
usr/local/share/virtualbox/nls/VirtualBox_sv.qm
usr/local/share/virtualbox/nls/VirtualBox_pt.qm
usr/local/share/virtualbox/nls/VirtualBox_fr.qm
usr/local/share/virtualbox/nls/VirtualBox_ca.qm
usr/local/share/virtualbox/nls/VirtualBox_nl.qm
usr/local/share/virtualbox/nls/qt_pl.qm
usr/local/share/virtualbox/nls/VirtualBox_uk.qm
usr/local/share/virtualbox/nls/qt_da.qm
usr/local/share/virtualbox/nls/qt_fr.qm
usr/local/share/virtualbox/nls/qt_sv.qm
usr/local/share/virtualbox/nls/qt_uk.qm
usr/local/share/virtualbox/nls/qt_eu.qm
usr/local/share/virtualbox/nls/VirtualBox_eu.qm
usr/local/share/virtualbox/nls/VirtualBox_pt_BR.qm
usr/local/share/virtualbox/nls/qt_es.qm
usr/local/share/virtualbox/nls/VirtualBox_sr.qm
usr/local/share/virtualbox/nls/qt_gl_ES.qm
usr/local/share/virtualbox/nls/qt_ro.qm
usr/local/share/virtualbox/nls/VirtualBox_en.qm
usr/local/share/virtualbox/nls/qt_tr.qm
usr/local/share/virtualbox/nls/VirtualBox_ar.qm
usr/local/share/virtualbox/nls/VirtualBox_km_KH.qm
usr/local/share/virtualbox/nls/qt_it.qm
usr/local/share/virtualbox/nls/qt_pt_BR.qm
usr/local/share/virtualbox/nls/VirtualBox_ja.qm
usr/local/share/virtualbox/nls/qt_id.qm
usr/local/share/virtualbox/nls/VirtualBox_lt.qm
usr/local/share/virtualbox/nls/VirtualBox_pl.qm
usr/local/share/virtualbox/nls/qt_ca_VA.qm
usr/local/share/virtualbox/nls/qt_ko.qm
usr/local/share/virtualbox/nls/VirtualBox_zh_TW.qm
usr/local/share/virtualbox/nls/qt_pt.qm
usr/local/share/virtualbox/nls/qt_ar.qm
usr/local/share/virtualbox/nls/VirtualBox_de.qm
usr/local/share/virtualbox/nls/VirtualBox_ru.qm
usr/local/share/virtualbox/nls/VirtualBox_id.qm
usr/local/share/applications/virtualbox-ose-beta.desktop
usr/local/share/virtualbox-ose-additions/xorg.conf.sample
usr/local/lib/VBoxOGLpassthroughspu.so
usr/local/lib/VBoxOGL.so
usr/local/lib/VBoxOGLpackspu.so
usr/local/lib/VBoxOGLerrorspu.so
usr/local/lib/security/pam_vbox.so
usr/local/lib/VBoxOGLfeedbackspu.so
usr/local/lib/modules/2.6.33.3-tinycore/kernel/misc/vboxdrv.ko
usr/local/lib/modules/2.6.33.3-tinycore/kernel/misc/vboxguest.ko
usr/local/lib/modules/2.6.33.3-tinycore/kernel/misc/vboxnetflt.ko
usr/local/lib/modules/2.6.33.3-tinycore/kernel/misc/vboxnetadp.ko
usr/local/lib/modules/2.6.33.3-tinycore/kernel/misc/vboxsf.ko
usr/local/lib/modules/2.6.33.3-tinycore/kernel/misc/vboxvideo.ko
usr/local/lib/X11/modules/drivers/vboxvideo.so
usr/local/lib/X11/modules/dri/vboxvideo_dri.so
usr/local/lib/X11/modules/input/vboxmouse.so
usr/local/lib/virtualbox/VBoxSDL.so
usr/local/lib/virtualbox/VBoxPython.so
usr/local/lib/virtualbox/VBoxREM64.so
usr/local/lib/virtualbox/VBoxOGLhostcrutil.so
usr/local/lib/virtualbox/VBoxDDR0.r0
usr/local/lib/virtualbox/VBoxXPCOMIPCD
usr/local/lib/virtualbox/VBoxREM.so
usr/local/lib/virtualbox/VirtualBox.so
usr/local/lib/virtualbox/VBoxAuthSimple.so
usr/local/lib/virtualbox/VBoxSVC
usr/local/lib/virtualbox/VBoxVMM.so
usr/local/lib/virtualbox/VBoxXPCOM.so
usr/local/lib/virtualbox/VBoxAuth.so
usr/local/lib/virtualbox/VBoxGuestControlSvc.so
usr/local/lib/virtualbox/VBoxTestOGL
usr/local/lib/virtualbox/VBoxDDU.so
usr/local/lib/virtualbox/VBoxKeyboard.so
usr/local/lib/virtualbox/VBoxSharedFolders.so
usr/local/lib/virtualbox/VBoxXPCOMC.so
usr/local/lib/virtualbox/VirtualBox
usr/local/lib/virtualbox/VBoxNetDHCP
usr/local/lib/virtualbox/VBoxDD.so
usr/local/lib/virtualbox/VBoxRT.so
usr/local/lib/virtualbox/VBoxSysInfo.sh
usr/local/lib/virtualbox/EfiThunk
usr/local/lib/virtualbox/VBoxHeadless
usr/local/lib/virtualbox/VBoxNetDHCP.so
usr/local/lib/virtualbox/VBoxNetAdpCtl
usr/local/lib/virtualbox/VBoxBFE
usr/local/lib/virtualbox/DBGCPlugInDiggers.so
usr/local/lib/virtualbox/VBoxEFI64.fd
usr/local/lib/virtualbox/VBoxEFI32.fd
usr/local/lib/virtualbox/VMMGC.gc
usr/local/lib/virtualbox/xpidl
usr/local/lib/virtualbox/VBoxDD2.so
usr/local/lib/virtualbox/VBoxDD2GC.gc
usr/local/lib/virtualbox/VMMR0.r0
usr/local/lib/virtualbox/VBoxSDL
usr/local/lib/virtualbox/VBoxREM32.so
usr/local/lib/virtualbox/VBoxManage
usr/local/lib/virtualbox/components/VBoxREM.so
usr/local/lib/virtualbox/components/VBoxC.so
usr/local/lib/virtualbox/components/VBoxVMM.so
usr/local/lib/virtualbox/components/VBoxXPCOM.so
usr/local/lib/virtualbox/components/VBoxDDU.so
usr/local/lib/virtualbox/components/VBoxSVCM.so
usr/local/lib/virtualbox/components/VBoxRT.so
usr/local/lib/virtualbox/components/VBoxXPCOMIPCC.so
usr/local/lib/virtualbox/components/VBoxXPCOMBase.xpt
usr/local/lib/virtualbox/components/VirtualBox_XPCOM.xpt
usr/local/lib/virtualbox/VBoxOGLhosterrorspu.so
usr/local/lib/virtualbox/VBoxDD2R0.r0
usr/local/lib/virtualbox/VBoxHeadless.so
usr/local/lib/virtualbox/VBoxDbg.so
usr/local/lib/virtualbox/VBoxGuestPropSvc.so
usr/local/lib/virtualbox/VBoxBFE.so
usr/local/lib/virtualbox/VBoxSharedCrOpenGL.so
usr/local/lib/virtualbox/VBoxOGLrenderspu.so
usr/local/lib/virtualbox/VBoxPython2_7.so
usr/local/lib/virtualbox/VBoxDDGC.gc
usr/local/lib/virtualbox/VBoxSharedClipboard.so
usr/local/lib/VBoxOGLarrayspu.so
usr/local/lib/VBoxOGLcrutil.so
usr/local/sbin/mount.vboxsf
usr/local/etc/xdg/autostart/vboxclient.desktop
usr/local/etc/vbox/vbox.cfg
usr/local/etc/init.d/vboxdrv
